# ops/com.kyrelldixon.kos-agent-restarter.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kyrelldixon.kos-agent-restarter</string>
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/kyrelldixon/projects/kos-agent/restart-kos-agent.sh</string>
    </array>
    <key>WatchPaths</key>
    <array>
        <string>/private/tmp/kos-agent-restart-trigger</string>
    </array>
    <key>StandardOutPath</key>
    <string>/Users/kyrelldixon/Library/Logs/kos-agent-restarter.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/kyrelldixon/Library/Logs/kos-agent-restarter.err</string>
</dict>
</plist>
